# layer/etc/init.d/dkl-boot
#!/sbin/openrc-run

description="Direktil boot operations."

depend()
{
    provide net
    before sysctl
}

start()
{
	ebegin "Running Direktil boot operations"
	/sbin/dkl init boot
	eend $? "Direktil boot operations failed"
}
